# decider-bam-slicer/osdc-decider.cron
#!/bin/bash

source ~/.bashrc

cd /glusterfs/netapp/homes1/BOCONNOR/gitroot/public-workflows/decider-bwa-pancancer
perl workflow_decider.pl --gnos-url https://gtrepo-osdc.annailabs.com --cluster-json my_cluster.json  --report osdc.log --upload-results
#perl generate_gnos_map.pl --output index.html --cluster-json my_cluster.json
#cp index.html /var/www/
#cp gtrepo*.log /var/www/

