# decider-bwa-pancancer/decider.cron
#!/bin/bash

source ~/.bashrc

cd /home/ubuntu/gitroot/public-workflows/decider-bwa-pancancer
perl workflow_decider.pl --gnos-url https://gtrepo-ebi.annailabs.com --cluster-json my_cluster.json  --report ebi.log --ignore-lane-count --upload-results
perl generate_gnos_map.pl --output index.html
cp index.html /var/www/

